# extra/cron.d/send-notifications
*/7 * * * *    root      ruby /var/inventario/script/runner -e production "NotificationsPool.send_notifications"
