# extra/cron.d/tracking-old-nodes
*/15 * * * *    root      ruby /var/inventario/script/runner -e production "Node.markOldEntries"
